# sql/crontab/mysql_diminishbans
#!/bin/sh

/usr/local/mysql/bin/mysql -h www.karchan.org -u karn --password=`cat /karchan/mud/sql/password.sql` -s mud <<END_OF_DATA

#
# Changeem
#
update bantable set days = days - 1 
	where days>0;
delete from bantable where days=0;

END_OF_DATA

